# contrib/user-local-install/files/.config/systemd/user/codex-desktop-update.timer
[Unit]
Description=Weekly update for Codex desktop Linux wrapper

[Timer]
OnCalendar=Sun 05:00
Persistent=true
Unit=codex-desktop-update.service

[Install]
WantedBy=timers.target
